# packaging/scripts/postinstall.sh
systemctl daemon-reload
systemctl enable cloudwatch-read-adapter.service
systemctl restart cloudwatch-read-adapter.service
