# packaging/scripts/preremove.sh
systemctl stop cloudwatch-read-adapter.service
systemctl disable cloudwatch-read-adapter.service
